# reboot_daily.timer
[Unit]
Description=Restart daily.

[Timer]
OnCalendar=*-*-* 03:03:00

[Install]
WantedBy=timers.target
