# reboot_daily.service
[Unit]
Description=Unit file for daily reboots

[Service]
Type=simple
ExecStart=/usr/bin/reboot
